# var/lib/dpkg/info/corosync.conffiles
/etc/corosync/corosync.conf.example.udpu
/etc/corosync/corosync.conf.example
/etc/corosync/corosync.conf
/etc/logrotate.d/corosync
/etc/default/corosync
/etc/init.d/corosync
