# var/lib/dpkg/info/drbdlinks.conffiles
/etc/init.d/drbdlinksclean
/etc/ha.d/resource.d/drbdlinks
/etc/drbdlinks.conf
