# var/lib/dpkg/info/drbd8-utils.list
/.
/etc
/etc/drbd.d
/etc/drbd.d/global_common.conf
/etc/bash_completion.d
/etc/bash_completion.d/drbdadm
/etc/ha.d
/etc/ha.d/resource.d
/etc/ha.d/resource.d/drbddisk
/etc/ha.d/resource.d/drbdupper
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/65-drbd.rules
/etc/init.d
/etc/init.d/drbd
/etc/drbd.conf
/etc/xen
/etc/xen/scripts
/etc/xen/scripts/block-drbd
/var
/var/lock
/var/lib
/var/lib/drbd
/sbin
/sbin/drbdsetup
/sbin/drbdadm
/sbin/drbdmeta
/usr
/usr/share
/usr/share/doc
/usr/share/doc/drbd8-utils
/usr/share/doc/drbd8-utils/copyright
/usr/share/doc/drbd8-utils/changelog.Debian.gz
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/drbd.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/drbdsetup.8.gz
/usr/share/man/man8/drbdmeta.8.gz
/usr/share/man/man8/drbdadm.8.gz
/usr/share/man/man8/drbddisk.8.gz
/usr/share/man/man8/drbd.8.gz
/usr/share/cluster
/usr/share/cluster/drbd.metadata
/usr/share/cluster/drbd.sh
/usr/sbin
/usr/sbin/drbd-overview
/usr/lib
/usr/lib/ocf
/usr/lib/ocf/resource.d
/usr/lib/ocf/resource.d/linbit
/usr/lib/ocf/resource.d/linbit/drbd
/usr/lib/drbd
/usr/lib/drbd/crm-fence-peer.sh
/usr/lib/drbd/snapshot-resync-target-lvm.sh
/usr/lib/drbd/notify.sh
/usr/lib/drbd/outdate-peer.sh
/usr/lib/drbd/notify-out-of-sync.sh
/usr/lib/drbd/notify-split-brain.sh
/usr/lib/drbd/notify-pri-lost.sh
/usr/lib/drbd/notify-emergency-shutdown.sh
/usr/lib/drbd/notify-pri-lost-after-sb.sh
/usr/lib/drbd/notify-pri-on-incon-degr.sh
/usr/lib/drbd/unsnapshot-resync-target-lvm.sh
/usr/lib/drbd/notify-emergency-reboot.sh
/usr/lib/drbd/notify-io-error.sh
/usr/lib/drbd/crm-unfence-peer.sh
