# var/lib/dpkg/info/drbdlinks.list
/.
/etc
/etc/init.d
/etc/init.d/drbdlinksclean
/etc/ha.d
/etc/ha.d/resource.d
/etc/ha.d/resource.d/drbdlinks
/etc/drbdlinks.conf
/usr
/usr/sbin
/usr/sbin/drbdlinks
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/drbdlinks.8.gz
/usr/share/doc
/usr/share/doc/drbdlinks
/usr/share/doc/drbdlinks/README
/usr/share/doc/drbdlinks/WHATSNEW
/usr/share/doc/drbdlinks/copyright
/usr/share/doc/drbdlinks/changelog.Debian.gz
