# var/lib/dpkg/info/kdump-tools.list
/.
/usr
/usr/sbin
/usr/sbin/kdump-config
/usr/share
/usr/share/doc
/usr/share/doc/kdump-tools
/usr/share/doc/kdump-tools/README.Debian
/usr/share/doc/kdump-tools/copyright
/usr/share/doc/kdump-tools/changelog.Debian.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/kdump-config.8.gz
/usr/share/man/man5
/usr/share/man/man5/kdump-tools.5.gz
/etc
/etc/default
/etc/default/kdump-tools
/etc/init.d
/etc/init.d/kdump-tools
