# var/lib/dpkg/info/kdump-tools.conffiles
/etc/default/kdump-tools
/etc/init.d/kdump-tools
