# var/lib/dpkg/info/kexec-tools.list
/.
/usr
/usr/lib
/usr/lib/kexec-tools
/usr/lib/kexec-tools/kexec_test.static
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/kexec.8.gz
/usr/share/man/man8/coldreboot.8.gz
/usr/share/man/man8/kdump.8.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/kexec-tools
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/0_kdump
/usr/share/doc
/usr/share/doc/kexec-tools
/usr/share/doc/kexec-tools/TODO
/usr/share/doc/kexec-tools/copyright
/usr/share/doc/kexec-tools/changelog.Debian.gz
/usr/share/doc/kexec-tools/News
/etc
/etc/init.d
/etc/init.d/kexec-load
/etc/init.d/kexec
/etc/init.d/kdump
/sbin
/sbin/kexec
/sbin/coldreboot
/sbin/kdump
