# var/lib/dpkg/info/mdadm.list
/.
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/64-md-raid.rules
/usr
/usr/share
/usr/share/man
/usr/share/man/man4
/usr/share/man/man4/md.4.gz
/usr/share/man/man8
/usr/share/man/man8/mdmon.8.gz
/usr/share/man/man8/mdadm.8.gz
/usr/share/man/man5
/usr/share/man/man5/mdadm.conf.5.gz
/usr/share/doc-base
/usr/share/doc-base/mdadm-raid5-vs-10
/usr/share/doc-base/mdadm-readme-recipes
/usr/share/doc-base/mdadm-md-txt
/usr/share/doc-base/mdadm-superblock-formats
/usr/share/doc-base/mdadm-jd-rebuild-raid
/usr/share/doc-base/mdadm-faq
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/mdadm
/usr/share/mdadm/checkarray
/usr/share/mdadm/mkconf
/usr/share/bug
/usr/share/bug/mdadm
/usr/share/bug/mdadm/script
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-premount
/usr/share/initramfs-tools/scripts/init-premount/mdadm
/usr/share/initramfs-tools/scripts/local-premount
/usr/share/initramfs-tools/scripts/local-premount/mdadm
/usr/share/initramfs-tools/scripts/mdadm-functions
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/mdadm
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_mdadm.py
/usr/share/doc
/usr/share/doc/mdadm
/usr/share/doc/mdadm/TODO.Debian
/usr/share/doc/mdadm/FAQ.gz
/usr/share/doc/mdadm/TODO.gz
/usr/share/doc/mdadm/md.txt.gz
/usr/share/doc/mdadm/copyright
/usr/share/doc/mdadm/NEWS.Debian.gz
/usr/share/doc/mdadm/README.Debian.gz
/usr/share/doc/mdadm/md_superblock_formats.txt.gz
/usr/share/doc/mdadm/examples
/usr/share/doc/mdadm/examples/mdadd.sh
/usr/share/doc/mdadm/examples/mdadm.conf-example
/usr/share/doc/mdadm/examples/syslog-events
/usr/share/doc/mdadm/changelog.Debian.gz
/usr/share/doc/mdadm/README.recipes.gz
/usr/share/doc/mdadm/RAID5_versus_RAID10.txt.gz
/usr/share/doc/mdadm/rebuilding-raid.html
/usr/share/doc/mdadm/README.checkarray
/etc
/etc/init.d
/etc/init.d/mdadm
/etc/cron.d
/etc/cron.d/mdadm
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/mdadm
/etc/logcheck/violations.d
/etc/logcheck/violations.d/mdadm
/etc/mdadm
/etc/cron.daily
/etc/cron.daily/mdadm
/sbin
/sbin/mdadm
/sbin/mdmon
/var
/var/lib
/var/lib/mdadm
