# var/lib/dpkg/info/mdadm.conffiles
/etc/init.d/mdadm
/etc/cron.d/mdadm
/etc/logcheck/ignore.d.server/mdadm
/etc/logcheck/violations.d/mdadm
/etc/cron.daily/mdadm
