# var/lib/dpkg/info/openipmi.conffiles
/etc/default/openipmi
/etc/init.d/openipmi
