# var/lib/dpkg/info/open-iscsi.conffiles
/etc/init.d/open-iscsi
/etc/iscsi/iscsid.conf
/etc/init/iscsi-network-interface.conf
