# var/lib/dpkg/info/openipmi.list
/.
/usr
/usr/bin
/usr/bin/ipmi_ui
/usr/bin/ipmilan
/usr/bin/openipmicmd
/usr/bin/openipmish
/usr/bin/rmcp_ping
/usr/bin/solterm
/usr/sbin
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/openipmish.1.gz
/usr/share/man/man1/rmcp_ping.1.gz
/usr/share/man/man1/solterm.1.gz
/usr/share/man/man1/ipmi_ui.1.gz
/usr/share/man/man1/openipmicmd.1.gz
/usr/share/man/man7
/usr/share/man/man7/ipmi_cmdlang.7.gz
/usr/share/man/man7/openipmi_conparms.7.gz
/usr/share/man/man8
/usr/share/man/man8/ipmilan.8.gz
/usr/share/doc
/usr/share/doc/openipmi
/usr/share/doc/openipmi/FAQ.gz
/usr/share/doc/openipmi/README.gz
/usr/share/doc/openipmi/README.Force
/usr/share/doc/openipmi/README.MotorolaMXP
/usr/share/doc/openipmi/TODO
/usr/share/doc/openipmi/copyright
/usr/share/doc/openipmi/changelog.Debian.gz
/etc
/etc/default
/etc/default/openipmi
/etc/init.d
/etc/init.d/openipmi
/usr/bin/ipmicmd
/usr/bin/ipmish
/usr/share/man/man1/ipmish.1.gz
/usr/share/man/man1/ipmicmd.1.gz
