# var/lib/dpkg/info/pacemaker.list
/.
/etc
/etc/init.d
/etc/init.d/pacemaker
/var
/var/lib
/var/lib/heartbeat
/var/lib/heartbeat/crm
/var/lib/heartbeat/lrm
/var/lib/heartbeat/pengine
/var/lib/pengine
/usr
/usr/share
/usr/share/doc
/usr/share/doc/pacemaker
/usr/share/doc/pacemaker/AUTHORS
/usr/share/doc/pacemaker/crm_fencing.txt.gz
/usr/share/doc/pacemaker/copyright
/usr/share/doc/pacemaker/README.hb2openais.gz
/usr/share/doc/pacemaker/acls.txt.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/pacemaker
/usr/share/pacemaker
/usr/share/pacemaker/templates
/usr/share/pacemaker/templates/virtual-ip
/usr/share/pacemaker/templates/gfs2-base
/usr/share/pacemaker/templates/clvm
/usr/share/pacemaker/templates/apache
/usr/share/pacemaker/templates/filesystem
/usr/share/pacemaker/templates/sbd
/usr/share/pacemaker/templates/gfs2
/usr/share/pacemaker/templates/ocfs2
/usr/share/pacemaker/resources-1.0.rng
/usr/share/pacemaker/constraints-1.1.rng
/usr/share/pacemaker/nvset.rng
/usr/share/pacemaker/constraints-1.0.rng
/usr/share/pacemaker/pacemaker-1.2.rng
/usr/share/pacemaker/crm.dtd
/usr/share/pacemaker/pacemaker-1.0.rng
/usr/share/pacemaker/resources-1.1.rng
/usr/share/pacemaker/upgrade06.xsl
/usr/share/pacemaker/acls-1.1.rng
/usr/share/pacemaker/resources-1.2.rng
/usr/share/pacemaker/constraints-1.2.rng
/usr/share/pacemaker/score.rng
/usr/share/pacemaker/pacemaker.rng
/usr/share/pacemaker/report.collector
/usr/share/pacemaker/pacemaker-1.1.rng
/usr/share/pacemaker/crm-transitional.dtd
/usr/share/pacemaker/acls-1.2.rng
/usr/share/pacemaker/report.common
/usr/share/pacemaker/rule.rng
/usr/share/pacemaker/versions.rng
/usr/share/pacemaker/crm.8.txt
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/crm_failcount.8.gz
/usr/share/man/man8/crm_uuid.8.gz
/usr/share/man/man8/crm_master.8.gz
/usr/share/man/man8/crm_node.8.gz
/usr/share/man/man8/crm_shadow.8.gz
/usr/share/man/man8/pacemakerd.8.gz
/usr/share/man/man8/crm_standby.8.gz
/usr/share/man/man8/fence_legacy.8.gz
/usr/share/man/man8/cibadmin.8.gz
/usr/share/man/man8/crm_ticket.8.gz
/usr/share/man/man8/iso8601.8.gz
/usr/share/man/man8/crm_mon.8.gz
/usr/share/man/man8/crm_report.8.gz
/usr/share/man/man8/crmadmin.8.gz
/usr/share/man/man8/crm_verify.8.gz
/usr/share/man/man8/attrd_updater.8.gz
/usr/share/man/man8/crm_simulate.8.gz
/usr/share/man/man8/stonith_admin.8.gz
/usr/share/man/man8/ptest.8.gz
/usr/share/man/man8/crm_diff.8.gz
/usr/share/man/man8/crm_resource.8.gz
/usr/share/man/man8/crm_attribute.8.gz
/usr/share/man/man7
/usr/share/man/man7/crmd.7.gz
/usr/share/man/man7/ocf_pacemaker_Dummy.7.gz
/usr/share/man/man7/ocf_pacemaker_SysInfo.7.gz
/usr/share/man/man7/ocf_pacemaker_o2cb.7.gz
/usr/share/man/man7/ocf_pacemaker_ping.7.gz
/usr/share/man/man7/ocf_pacemaker_ClusterMon.7.gz
/usr/share/man/man7/ocf_pacemaker_controld.7.gz
/usr/share/man/man7/pengine.7.gz
/usr/share/man/man7/ocf_pacemaker_HealthCPU.7.gz
/usr/share/man/man7/ocf_pacemaker_pingd.7.gz
/usr/share/man/man7/ocf_pacemaker_SystemHealth.7.gz
/usr/share/man/man7/ocf_pacemaker_HealthSMART.7.gz
/usr/share/man/man7/stonithd.7.gz
/usr/share/man/man7/ocf_pacemaker_Stateful.7.gz
/usr/share/pyshared
/usr/share/pyshared/crm
/usr/share/pyshared/crm/term.py
/usr/share/pyshared/crm/ui.py
/usr/share/pyshared/crm/completion.py
/usr/share/pyshared/crm/template.py
/usr/share/pyshared/crm/userprefs.py
/usr/share/pyshared/crm/main.py
/usr/share/pyshared/crm/help.py
/usr/share/pyshared/crm/levels.py
/usr/share/pyshared/crm/ra.py
/usr/share/pyshared/crm/idmgmt.py
/usr/share/pyshared/crm/cache.py
/usr/share/pyshared/crm/clidisplay.py
/usr/share/pyshared/crm/vars.py
/usr/share/pyshared/crm/__init__.py
/usr/share/pyshared/crm/singletonmixin.py
/usr/share/pyshared/crm/msg.py
/usr/share/pyshared/crm/xmlutil.py
/usr/share/pyshared/crm/parse.py
/usr/share/pyshared/crm/utils.py
/usr/share/pyshared/crm/cibconfig.py
/usr/share/pyshared/crm/cliformat.py
/usr/share/pyshared/crm/cibstatus.py
/usr/share/pyshared/cts
/usr/share/pyshared/cts/CTSvars.py
/usr/share/pyshared/cts/CM_ais.py
/usr/share/pyshared/cts/CTSscenarios.py
/usr/share/pyshared/cts/CTS.py
/usr/share/pyshared/cts/CTSaudits.py
/usr/share/pyshared/cts/CM_lha.py
/usr/share/pyshared/cts/CIB.py
/usr/share/pyshared/cts/__init__.py
/usr/share/pyshared/cts/CTStests.py
/usr/share/snmp
/usr/share/snmp/mibs
/usr/share/snmp/mibs/PCMK-MIB.txt
/usr/sbin
/usr/sbin/crm_verify
/usr/sbin/fence_legacy
/usr/sbin/crm_report
/usr/sbin/cibpipe
/usr/sbin/crm_mon
/usr/sbin/crm_uuid
/usr/sbin/crm_failcount
/usr/sbin/crm
/usr/sbin/crm_node
/usr/sbin/crmadmin
/usr/sbin/cibadmin
/usr/sbin/crm_ticket
/usr/sbin/crm_simulate
/usr/sbin/crm_master
/usr/sbin/crm_attribute
/usr/sbin/attrd_updater
/usr/sbin/stonith_admin
/usr/sbin/crm_shadow
/usr/sbin/crm_resource
/usr/sbin/pacemakerd
/usr/sbin/crm_diff
/usr/sbin/crm_standby
/usr/sbin/ptest
/usr/sbin/iso8601
/usr/lib
/usr/lib/heartbeat
/usr/lib/heartbeat/plugins
/usr/lib/heartbeat/plugins/RAExec
/usr/lib/heartbeat/plugins/RAExec/stonith.a
/usr/lib/heartbeat/plugins/RAExec/stonith.so
/usr/lib/heartbeat/plugins/RAExec/stonith.la
/usr/lib/heartbeat/stonithd
/usr/lib/heartbeat/attrd
/usr/lib/heartbeat/crmd
/usr/lib/heartbeat/cibmon
/usr/lib/heartbeat/hb2openais.sh
/usr/lib/heartbeat/atest
/usr/lib/heartbeat/crm_primitive.py
/usr/lib/heartbeat/stonith-test
/usr/lib/heartbeat/pengine
/usr/lib/heartbeat/hb2openais-helper.py
/usr/lib/heartbeat/cib
/usr/lib/lcrso
/usr/lib/lcrso/pacemaker.lcrso
/usr/lib/ocf
/usr/lib/ocf/resource.d
/usr/lib/ocf/resource.d/pacemaker
/usr/lib/ocf/resource.d/pacemaker/controld
/usr/lib/ocf/resource.d/pacemaker/HealthCPU
/usr/lib/ocf/resource.d/pacemaker/pingd
/usr/lib/ocf/resource.d/pacemaker/HealthSMART
/usr/lib/ocf/resource.d/pacemaker/SystemHealth
/usr/lib/ocf/resource.d/pacemaker/SysInfo
/usr/lib/ocf/resource.d/pacemaker/Dummy
/usr/lib/ocf/resource.d/pacemaker/ClusterMon
/usr/lib/ocf/resource.d/pacemaker/Stateful
/usr/lib/ocf/resource.d/pacemaker/ping
/usr/lib/ocf/resource.d/pacemaker/o2cb
/usr/lib/python2.7
/usr/lib/python2.7/dist-packages
/usr/lib/python2.7/dist-packages/crm
/usr/lib/python2.7/dist-packages/cts
/usr/share/doc/pacemaker/changelog.Debian.gz
/usr/lib/python2.7/dist-packages/crm/term.py
/usr/lib/python2.7/dist-packages/crm/ui.py
/usr/lib/python2.7/dist-packages/crm/completion.py
/usr/lib/python2.7/dist-packages/crm/template.py
/usr/lib/python2.7/dist-packages/crm/userprefs.py
/usr/lib/python2.7/dist-packages/crm/main.py
/usr/lib/python2.7/dist-packages/crm/help.py
/usr/lib/python2.7/dist-packages/crm/levels.py
/usr/lib/python2.7/dist-packages/crm/ra.py
/usr/lib/python2.7/dist-packages/crm/idmgmt.py
/usr/lib/python2.7/dist-packages/crm/cache.py
/usr/lib/python2.7/dist-packages/crm/clidisplay.py
/usr/lib/python2.7/dist-packages/crm/vars.py
/usr/lib/python2.7/dist-packages/crm/__init__.py
/usr/lib/python2.7/dist-packages/crm/singletonmixin.py
/usr/lib/python2.7/dist-packages/crm/msg.py
/usr/lib/python2.7/dist-packages/crm/xmlutil.py
/usr/lib/python2.7/dist-packages/crm/parse.py
/usr/lib/python2.7/dist-packages/crm/utils.py
/usr/lib/python2.7/dist-packages/crm/cibconfig.py
/usr/lib/python2.7/dist-packages/crm/cliformat.py
/usr/lib/python2.7/dist-packages/crm/cibstatus.py
/usr/lib/python2.7/dist-packages/cts/CTSvars.py
/usr/lib/python2.7/dist-packages/cts/CM_ais.py
/usr/lib/python2.7/dist-packages/cts/CTSscenarios.py
/usr/lib/python2.7/dist-packages/cts/CTS.py
/usr/lib/python2.7/dist-packages/cts/CTSaudits.py
/usr/lib/python2.7/dist-packages/cts/CM_lha.py
/usr/lib/python2.7/dist-packages/cts/CIB.py
/usr/lib/python2.7/dist-packages/cts/__init__.py
/usr/lib/python2.7/dist-packages/cts/CTStests.py
