# var/lib/dpkg/info/postgresql-common.conffiles
/etc/init.d/postgresql
/etc/sysctl.d/30-postgresql-shm.conf
/etc/logrotate.d/postgresql-common
