# var/lib/dpkg/info/postgresql-common.list
/.
/var
/var/log
/var/log/postgresql
/var/lib
/var/lib/postgresql
/usr
/usr/sbin
/usr/sbin/pg_updatedicts
/usr/bin
/usr/bin/pg_config
/usr/bin/pg_upgradecluster
/usr/bin/pg_virtualenv
/usr/bin/pg_dropcluster
/usr/bin/pg_createcluster
/usr/bin/pg_ctlcluster
/usr/bin/pg_lsclusters
/usr/share
/usr/share/postgresql-common
/usr/share/postgresql-common/pg_checksystem
/usr/share/postgresql-common/run-upgrade-scripts
/usr/share/postgresql-common/init.d-functions
/usr/share/postgresql-common/upgrade-scripts
/usr/share/postgresql-common/upgrade-scripts/SPECIFICATION
/usr/share/postgresql-common/t
/usr/share/postgresql-common/t/130_nonroot_admin.t
/usr/share/postgresql-common/t/030_errors.t
/usr/share/postgresql-common/t/TestLib.pm
/usr/share/postgresql-common/t/140_pg_config.t
/usr/share/postgresql-common/t/041_upgrade_custompaths.t
/usr/share/postgresql-common/t/042_upgrade_tablespaces.t
/usr/share/postgresql-common/t/052_upgrade_encodings.t
/usr/share/postgresql-common/t/150_tsearch_stemming.t
/usr/share/postgresql-common/t/template
/usr/share/postgresql-common/t/020_create_sql_remove.t
/usr/share/postgresql-common/t/001_packages.t
/usr/share/postgresql-common/t/003_package_checks.t
/usr/share/postgresql-common/t/010_defaultport_cluster.t
/usr/share/postgresql-common/t/031_errors_disk_full.t
/usr/share/postgresql-common/t/110_integrate_cluster.t
/usr/share/postgresql-common/t/100_upgrade_scripts.t
/usr/share/postgresql-common/t/180_ecpg.t
/usr/share/postgresql-common/t/170_extensions.t
/usr/share/postgresql-common/t/050_encodings.t
/usr/share/postgresql-common/t/085_pg_ctl.conf.t
/usr/share/postgresql-common/t/080_start.conf.t
/usr/share/postgresql-common/t/160_alternate_confroot.t
/usr/share/postgresql-common/t/060_obsolete_confparams.t
/usr/share/postgresql-common/t/040_upgrade.t
/usr/share/postgresql-common/t/051_inconsistent_encoding_upgrade.t
/usr/share/postgresql-common/t/070_non_postgres_clusters.t
/usr/share/postgresql-common/t/005_PgCommon.t
/usr/share/postgresql-common/t/090_multicluster.t
/usr/share/postgresql-common/t/002_existing_clusters.t
/usr/share/postgresql-common/t/120_pg_upgradecluster_scripts.t
/usr/share/postgresql-common/testsuite
/usr/share/doc
/usr/share/doc/postgresql-common
/usr/share/doc/postgresql-common/copyright
/usr/share/doc/postgresql-common/architecture.html
/usr/share/doc/postgresql-common/changelog.gz
/usr/share/doc/postgresql-common/README.Devel
/usr/share/doc/postgresql-common/README.Debian.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/pg_virtualenv.1.gz
/usr/share/man/man1/pg_lsclusters.1.gz
/usr/share/man/man8
/usr/share/man/man8/pg_dropcluster.8.gz
/usr/share/man/man8/pg_updatedicts.8.gz
/usr/share/man/man8/pg_upgradecluster.8.gz
/usr/share/man/man8/pg_ctlcluster.8.gz
/usr/share/man/man8/pg_createcluster.8.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/postgresql-common
/etc
/etc/postgresql-common
/etc/postgresql-common/pg_upgradecluster.d
/etc/init.d
/etc/init.d/postgresql
/etc/sysctl.d
/etc/sysctl.d/30-postgresql-shm.conf
/etc/logrotate.d
/etc/logrotate.d/postgresql-common
