# var/lib/dpkg/info/rabbitmq-server.conffiles
/etc/init.d/rabbitmq-server
/etc/rabbitmq/rabbitmq-env.conf
/etc/logrotate.d/rabbitmq-server
