# var/lib/dpkg/info/rabbitmq-erlang-client.postinst
#!/bin/sh -e

/usr/lib/rabbitmq/bin/rabbitmq-plugins enable amqp_client
invoke-rc.d rabbitmq-server restart


