# var/lib/dpkg/info/tomcat7.list
/.
/usr
/usr/share
/usr/share/tomcat7-root
/usr/share/tomcat7-root/default_root
/usr/share/tomcat7-root/default_root/META-INF
/usr/share/tomcat7-root/default_root/META-INF/context.xml
/usr/share/tomcat7-root/default_root/index.html
/usr/share/tomcat7
/usr/share/tomcat7/defaults.template
/usr/share/tomcat7/defaults.md5sum
/usr/share/doc
/usr/share/doc/tomcat7
/usr/share/doc/tomcat7/copyright
/var
/var/log
/var/log/tomcat7
/var/cache
/var/cache/tomcat7
/var/lib
/var/lib/tomcat7
/var/lib/tomcat7/common
/var/lib/tomcat7/common/classes
/var/lib/tomcat7/server
/var/lib/tomcat7/server/classes
/var/lib/tomcat7/webapps
/var/lib/tomcat7/shared
/var/lib/tomcat7/shared/classes
/etc
/etc/init.d
/etc/init.d/tomcat7
/etc/tomcat7
/etc/tomcat7/logging.properties
/etc/tomcat7/server.xml
/etc/tomcat7/web.xml
/etc/tomcat7/catalina.properties
/etc/tomcat7/context.xml
/etc/tomcat7/Catalina
/etc/tomcat7/Catalina/localhost
/etc/tomcat7/tomcat-users.xml
/etc/tomcat7/policy.d
/etc/tomcat7/policy.d/01system.policy
/etc/tomcat7/policy.d/02debian.policy
/etc/tomcat7/policy.d/03catalina.policy
/etc/tomcat7/policy.d/04webapps.policy
/etc/tomcat7/policy.d/50local.policy
/etc/logrotate.d
/etc/logrotate.d/tomcat7
/etc/cron.daily
/etc/cron.daily/tomcat7
/usr/share/doc/tomcat7/README.Debian.gz
/usr/share/doc/tomcat7/changelog.Debian.gz
/var/lib/tomcat7/work
/var/lib/tomcat7/logs
/var/lib/tomcat7/conf
