# var/lib/dpkg/info/watchdog.conffiles
/etc/watchdog.conf
/etc/init.d/wd_keepalive
/etc/init.d/watchdog
