# var/lib/dpkg/info/watchdog.list
/.
/etc
/etc/watchdog.conf
/etc/init.d
/etc/init.d/wd_keepalive
/etc/init.d/watchdog
/var
/var/log
/var/log/watchdog
/usr
/usr/share
/usr/share/doc
/usr/share/doc/watchdog
/usr/share/doc/watchdog/watch_err.h
/usr/share/doc/watchdog/IAFA-PACKAGE
/usr/share/doc/watchdog/examples
/usr/share/doc/watchdog/examples/another-chance.sh
/usr/share/doc/watchdog/examples/uptime.sh
/usr/share/doc/watchdog/examples/systemcheck.sh
/usr/share/doc/watchdog/examples/README
/usr/share/doc/watchdog/examples/repair.sh
/usr/share/doc/watchdog/examples/dbcheck.sh
/usr/share/doc/watchdog/watchdog.lsm
/usr/share/doc/watchdog/copyright
/usr/share/doc/watchdog/README.watchdog.ipmi.gz
/usr/share/doc/watchdog/README
/usr/share/doc/watchdog/changelog.Debian.gz
/usr/share/doc/watchdog/README.Debian
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/watchdog.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/wd_keepalive.8.gz
/usr/share/man/man8/watchdog.8.gz
/usr/share/man/man8/wd_identify.8.gz
/usr/sbin
/usr/sbin/wd_identify
/usr/sbin/wd_keepalive
/usr/sbin/watchdog
